# docker/vnc/bg.conf
[program:feh]
command: feh --bg-fill /app/SBEL.png
autostart=true
autorestart=false
startretries=2
